# .agent/cron/digital-garden.cron
# 数字花园孵化器 - 定时任务配置
# 4个独立时间点，按顺序执行
# 安装: crontab /opt/digital-garden-incubation-zone/.agent/cron/digital-garden.cron

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# 部署密码 (请修改为实际密码)
SERVER_PASSWORD=YOUR_PASSWORD_HERE

# 05:00 - 阶段1: 灵感捕获
0 5 * * * /opt/digital-garden-incubation-zone/.agent/scripts/step1-inspiration.sh >> /opt/digital-garden-incubation-zone/.agent/logs/cron.log 2>&1

# 05:30 - 阶段2: 创意孵化
30 5 * * * /opt/digital-garden-incubation-zone/.agent/scripts/step2-incubate.sh >> /opt/digital-garden-incubation-zone/.agent/logs/cron.log 2>&1

# 06:00 - 阶段3: 委派研发
0 6 * * * /opt/digital-garden-incubation-zone/.agent/scripts/step3-develop.sh >> /opt/digital-garden-incubation-zone/.agent/logs/cron.log 2>&1

# 08:00 - 阶段4: 部署
0 8 * * * /opt/digital-garden-incubation-zone/.agent/scripts/step4-deploy.sh >> /opt/digital-garden-incubation-zone/.agent/logs/cron.log 2>&1
